# packaging/systemd/hermes-zero-token.service
# hermes zero-token OAuth proxy (systemd *user* service)
#
# A standalone OpenAI-/Anthropic-compatible endpoint that authenticates to
# Anthropic with the Claude CLI *account* token (claude /login / setup-token)
# instead of a paid API key. Point hermes-agent's `provider: custom` backend at
# http://127.0.0.1:3031/v1 (see docs/zero-token-claude-oauth.md).
#
# Install:
#   cp packaging/systemd/hermes-zero-token.service ~/.config/systemd/user/
#   # edit WorkingDirectory / venv path below to match your checkout
#   systemctl --user daemon-reload
#   systemctl --user enable --now hermes-zero-token
#
# Requires: the `zero-token` extra (aiohttp). Either install it into the hermes
# venv (`uv pip install -e '.[zero-token]'`) or run from any venv that has
# aiohttp + httpx.

[Unit]
Description=hermes zero-token OAuth proxy (Claude account token -> OpenAI/Anthropic API)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=%h/hermes-agent
ExecStart=%h/hermes-agent/.venv/bin/python -m zero_token

# --- configuration (edit to taste) ---------------------------------------
# Bind address/port. 3031 matches the legacy claude_code_proxy default so an
# existing hermes `model.base_url: http://127.0.0.1:3031/v1` keeps working.
Environment=ZT_HOST=127.0.0.1
Environment=ZT_PORT=3031
# Default model when the client doesn't specify one.
Environment=ZT_DEFAULT_MODEL=claude-opus-4-8
# Shared-secret bearer token. MUST match hermes `model.api_key`. Generate with
# `openssl rand -hex 32`. Without it the proxy is unauthenticated (it can spend
# your Claude subscription — always set one).
Environment=CLAUDE_PROXY_TOKEN=CHANGE_ME_openssl_rand_hex_32
# Per-request upstream timeout (seconds). Each request is a single model turn
# (hermes drives its own tool loop), so this can be modest.
Environment=ZT_REQUEST_TIMEOUT=300
Environment=LOG_LEVEL=INFO

# The proxy reads ~/.claude/.credentials.json and refreshes it in place; no
# ANTHROPIC_API_KEY is needed. To use a static setup-token instead, set:
#   Environment=ANTHROPIC_TOKEN=sk-ant-oat01-...
#
# Multi-account failover: when the primary account hits "out of extra usage",
# the proxy rotates to the next account. Configure an ordered JSON array
# (primary first). Providers: anthropic (Claude sub, OAuth account token), kimi
# (Kimi Code sub, a plain API key from https://www.kimi.com/code/en), generic.
# claude2 here uses a static setup-token (claude setup-token); alternatively
# point it at a second credentials file via "credentials_path".
#
# PREFERRED: keep the array (with its tokens) in a 0600 file and reference it,
# so secrets stay off the systemd env / `systemctl show`:
#   Environment=ZT_ACCOUNTS_FILE=%h/.hermes/zt-accounts.json
# where that file contains, e.g.:
#   [{"name":"claude1","provider":"anthropic"},
#    {"name":"claude2","provider":"anthropic","token":"sk-ant-oat01-..."},
#    {"name":"kimi","provider":"kimi","token":"<KIMI_API_KEY>","model":"kimi-code"}]
#
# Inline alternative (quoting is fragile; secrets visible via systemctl show):
#   Environment=ZT_ACCOUNTS_JSON=[{"name":"claude1","provider":"anthropic"}, ...]

Restart=always
RestartSec=5
# Ignore SIGHUP (also enforced in-process) so a terminal hangup can't kill it.
StandardInput=null
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
